# cmd/DEBIAN/postinst
#!/bin/sh

set -e

systemctl daemon-reload
systemctl enable proxyd.service
systemctl start proxyd.service

echo "proxyd service installed and started."
